# dogwatch.service
[Unit]
Description=DogWatch Dog Detection Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=sf27
WorkingDirectory=/home/sf27/dogdetector9000
ExecStart=/usr/bin/python3 /home/sf27/dogdetector9000/dog_watch.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
